# support/WebTools/Library/LaunchDeamons/com.github.winnebago.plistmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.winnebago.plistmonitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/killall -USR1</string>
		<string>python</string>
	</array>
	<key>QueueDirectories</key>
	<array/>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>WatchPaths</key>
	<array>
		<string>/Shared Items/MigrationLogs/.usersync/</string>
	</array>
</dict>
</plist>
